# apc_utility_auto_off.timer
[Unit]
Description=Power off APC PDU outlet after set time since boot

[Timer]
OnBootSec=3h

[Install]
WantedBy=timers.target
